# scripts/com.supernote-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- Unique identifier for the launchd job -->
    <key>Label</key>
    <string>com.supernote-sync</string>

    <!-- Command to run: supernote-sync watch using the venv binary -->
    <key>ProgramArguments</key>
    <array>
      <string>/Users/YOUR_USER/Documents/GitHub/supernote-to-obsidian/.venv/bin/supernote-sync</string>
      <string>--config</string>
      <string>/Users/YOUR_USER/Documents/GitHub/supernote-to-obsidian/config.yaml</string>
      <string>watch</string>
    </array>

    <!-- Start automatically when the user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/.supernote-sync/logs/supernote-sync.stdout.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/.supernote-sync/logs/supernote-sync.stderr.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USER/Documents/GitHub/supernote-to-obsidian</string>
  </dict>
</plist>
